# systemd/robertapp-backup.timer
[Unit]
Description=Run RobertApp SQLite backup daily

[Timer]
OnCalendar=*-*-* 03:15:00
Persistent=true
Unit=robertapp-backup.service

[Install]
WantedBy=timers.target
