# systemd/robertapp-backup.service
[Unit]
Description=RobertApp SQLite backup
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
Environment=APP_DIR=/opt/RobertApp
Environment=BACKUP_DIR=/opt/backups/robertapp
Environment=RETENTION_DAYS=30
ExecStart=/opt/RobertApp/scripts/backup-sqlite.sh
User=robert
Group=robert
